# spec/fixtures/plugins/api/supervisord.conf
[supervisord]
nodaemon=true
logfile=/api/supervisord.log;
pidfile=/var/run/supervisord.pid;
loglevel=debug

[program:slapin]
command=rackup -E production -o 0.0.0.0 -p 4700
directory=/api
autorestart=true
redirect_stderr=true
redirect_stdout=true
